# goaccess/rootfs/etc/services.d/goaccess/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start the GoAccess service
# s6-overlay docs: https://github.com/just-containers/s6-overlay
# ==============================================================================

# shellcheck shell=bash

bashio::log.info "Starting GoAccess..."

if [ -f /config/goaccess.conf ]; then
    exec goaccess --config-file=/config/goaccess.conf
else
    exec goaccess --config-file=/data/goaccess.conf
fi
